# orchestration/cube_pull_r2.service
# cube_pull_r2.service — archive-box side of the R2 cube pipeline as a systemd user unit
# (setup in orchestration/README.md). The script resolves rclone (~/bin fallback, RCLONE=
# override) and its creds (~/.config/edm-r2.env) itself; only PULL_DEST and the repo path
# are machine policy — override via `systemctl --user edit cube_pull_r2`.
[Unit]
Description=EDM R2 cube puller (download, sha256-verify, free bucket slots)
After=network-online.target
Wants=network-online.target

[Service]
Environment=PULL_DEST=/storage/pool/smc/edm_cubes
ExecStart=/usr/bin/env bash %h/.julia/dev/ElectronDynamicsModels/orchestration/cube_pull_r2.sh
Restart=always
RestartSec=60
StandardOutput=append:%h/cube_pull_r2.log
StandardError=append:%h/cube_pull_r2.log

[Install]
WantedBy=default.target
